# supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=/root/server
autorestart=true
